# assets/etc/supervisor/conf.d/dockerd.ini
[program:docker]
priority=10
command=/bin/bash -c "source /etc/profile.d/docker-environment.sh; /usr/local/bin/wrapdocker"
stopasgroup=true
startsecs=0
exitcodes=0,1
