# assets/etc/supervisor/conf.d/dockerd.ini
[program:docker]
priority=10
command=/usr/local/bin/wrapdocker
startsecs=0
exitcodes=0,1
