# assets/etc/supervisor/conf.d/jenkins.ini
[program:jenkins]
priority=30
user=jenkins
command=/bin/tini -s -- /usr/local/bin/jenkins.sh
