# assets/etc/supervisor/conf.d/jenkins.ini
[program:jenkins]
priority=30
user=jenkins
command=/bin/bash -c "source /etc/profile.d/docker-environment.sh; /bin/tini -s -- /usr/local/bin/jenkins.sh"
stopasgroup=true
